# supervisord.conf
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]               ; inet (TCP) server disabled by default
port=*:15000                     ; (ip_address:port specifier, *:port for all iface)
username=%(ENV_SUPERVISOR_USER)s ; (default is no username (open server))
password=%(ENV_SUPERVISOR_PWD)s  ; (default is no password (open server))


[supervisord]
directory=/var/
environment=DJANGO_SETTINGS_MODULE="sir.config.settings"
logfile=/logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_backups=1
logfile_maxbytes=512Kb
loglevel=info                ; (log level;default info; others: debug,warn,trace)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
nodaemon=true               ; (start in foreground if true;default false)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

[program:worker]
command=sentry run worker
autostart=%(ENV_START_WORKER)s
autorestart=false
numprocs=1
redirect_stderr=true
logfile_maxbytes=512Kb
priority=700
startretries=3
startsecs=10
stopwaitsecs=600
killasgroup=true
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.


[program:cron]
command=sentry run cron
autostart=%(ENV_START_CRON)s
autorestart=false
numprocs=1
redirect_stderr=true
#stdout_logfile=celery-beat.log
logfile_maxbytes=512Kb
priority=700
startretries=3
startsecs=10
stopwaitsecs=600
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.


[program:web]
command=sentry run web -w 4 --noinput
autostart=%(ENV_START_WEB)s
numprocs=1
autorestart=false
redirect_stderr=true
#stdout_logfile=flower.log
priority=600
